# root/etc/services.d/chatbot/run
#!/usr/bin/with-contenv bash

cd /app/chatbot

exec \
  s6-setuidgid abc \
    /usr/bin/java -jar \
      chatbot.jar