# root/etc/services.d/chatbot/run
#!/usr/bin/with-contenv bash

cd /app/chatbot

exec \
  s6-setuidgid grfd \
    /usr/bin/java -Dlog4j2.formatMsgNoLookups=true -jar \
      chatbot.jar